# test.crontab

# test ranges 
0-23,24,25,26,40-45,50-59 * * JAN-MAR,MAY-JUN,AUG-OCT,DEC SAT,2,SUN someCommand

# test wildcards
* * * * * some Command

# test @whatever
@reboot some Command
@yearly some Command
@monthly some Command
@weekly some Command
@daily some Command
@hourly some Command

# dom starts by 1 --> test that !
* * 3 * * some Command

# stepwise stuff
*/6 */5 */4 */3 */2 some Command
*/50 * */31 */12 */1 some Command
