# test.crontab

# test ranges 
0-23,24,25,26,40-45,50-59 * * JAN-MAR,MAY-JUN,AUG-OCT,DEC SAT,2,SUN someCommand

# test wildcards
* * * * * some Command

# test @whatever
@reboot some Command
@yearly some Command
@monthly some Command
@weekly some Command
@daily some Command
@hourly some Command

# test stepwise stuff
2-40/2 * * * * some Command
0-20/2,30-40/2 * * * * some Command
*/6 */5 */4 */3 */2 some Command
